# deploy/agent-company-product-runner.service
[Unit]
Description=Agent Company Product Engineer execution runner
After=local-fs.target agent-company-worker.service

[Service]
Type=simple
WorkingDirectory=%h/agent-company
Environment="HOME=%h"
Environment="PATH=%h/.npm-global/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/env python3.11 -m agent_company.cli runner-run --executor-id product-engineer-runner --owner "Product Engineer" --capability product --capability engineering --poll-seconds 5
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
